# ozon-positions.timer
[Unit]
Description=Run Ozon positions parser every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
